# ansible-eth-monitoring/roles/version_check/templates/version-check.timer.j2
[Unit]
Description=Ethereum Client Version Check Timer

[Timer]
OnUnitActiveSec={{ version_check_interval }}
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target
